# Systemd/run-etherpad.service
[Unit]
description="A service that runs etherpad when the server is launched"

[Service]
User=root
WorkingDirectory=
